# ops/systemd/predictfun-mm-monitor.service
[Unit]
Description=Predict.fun MM monitor web service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/predictfun-mm/predictfun-mm.env
ExecStart=/usr/bin/env bash -lc 'cd "$PREDICT_MM_WORKDIR" && npm run monitor-web -- --host=127.0.0.1 --port=8787'
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
