# ops/systemd/predictfun-mm-shadow.service
[Unit]
Description=Predict.fun MM shadow runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/predictfun-mm/predictfun-mm.env
ExecStart=/usr/bin/env bash -lc 'cd "$PREDICT_MM_WORKDIR" && npm run shadow'
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
